# service/chappie.service
[Unit]
Description = chappie daemon

[Service]
Restart = no
Type = simple
ExecStart =
User =

[Install]
WantedBy = multi-user.target

